# deploy/supervisord.conf
#********************************************************************
# 项目名称：JeriBlog
# 文件名称：supervisord.conf
# 创建时间：2026-03-31 11:50:09
#
# 系统用户：Administrator
# 作　　者：Jerion
# 联系邮箱：416685476@qq.com
# 功能描述：supervisord 多进程管理配置，统一管理 Go 服务、Nuxt SSR 与 Nginx
#********************************************************************

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

# ── Go 后端服务 ────────────────────────────────────────────────────
# 不设置 environment，继承容器全部环境变量（DB_HOST / JWT_SECRET 等）
[program:server]
command=/app/blog-backend
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/server.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

# ── Nuxt 3 SSR 博客前台 ────────────────────────────────────────────
# 显式传入运行时必要变量；NUXT_PUBLIC_API_URL 从容器环境透传
[program:blog]
command=node /app/blog/.output/server/index.mjs
directory=/app
autostart=true
autorestart=true
environment=NODE_ENV="production",PORT="3000",HOST="127.0.0.1",NUXT_PUBLIC_API_URL="%(ENV_NUXT_PUBLIC_API_URL)s"
stdout_logfile=/app/logs/blog.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/blog.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3

# ── Nginx 反向代理 ─────────────────────────────────────────────────
# 日志直接输出到 docker logs
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
